# distrobox-update.service
[Unit]
Description=Distrobox Upgrade Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/distrobox-upgrade --all
User=your_username

[Install]
WantedBy=multi-user.target
